# com.neolight.guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neolight.guardian</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>~/neolight/neo_light_fix.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>~/neolight/logs/guardian_launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>~/neolight/logs/guardian_launchd_error.log</string>
    
    <key>WorkingDirectory</key>
    <string>~/neolight</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:~/neolight/venv/bin</string>
        <key>PYTHONPATH</key>
        <string>~/neolight</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    
    <key>StartInterval</key>
    <integer>300</integer>
    
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>

